# scripts/deploy/luck-agent.service
[Unit]
Description=Luck Agent tick — Pepe social voice for @0xleonardo
After=network.target stack-xi-pepe.service

[Service]
Type=oneshot
WorkingDirectory=/var/www/pepe-buildingculture
EnvironmentFile=/var/www/pepe-buildingculture/.env
ExecStart=/usr/bin/node scripts/luck-agent-tick.mjs
StandardOutput=append:/var/log/pepe-buildingculture/luck-agent.log
StandardError=append:/var/log/pepe-buildingculture/luck-agent.log

[Install]
WantedBy=multi-user.target
